# tools/launchd/com.thatcube.reclaim-disk.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thatcube.reclaim-disk</string>

    <!-- Runs the shared daily build-cache reclaim for Plozz/Mozz/Twozz.
         Installed by tools/install-reclaim-agent.sh (which rewrites the
         ProgramArguments path to this machine's checkout). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/brandon/Development/Plozz/tools/reclaim-disk.sh</string>
        <string>--days</string>
        <string>4</string>
    </array>

    <!-- Daily at 03:30. If the Mac is asleep at that time, launchd runs it
         at the next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <!-- Low priority; this is background maintenance. -->
    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>10</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/brandon/Library/Logs/reclaim-disk.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/brandon/Library/Logs/reclaim-disk.log</string>
</dict>
</plist>
